# .config/systemd/user/duplicacy-backup-tom.service
[Unit]
Description=Perform duplicacy backup

[Service]
ExecStart=/home/tom/bin/duplicacy-backup.sh
